# recipes-devtools/python/files/python-barbican/openstack-barbican-worker.service
[Unit]
Description=Openstack Barbican worker daemon
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/bin/barbican-worker
User=barbican
Group=barbican
Restart=on-failure

[Install]
WantedBy=multi-user.target
